# NODIRSTUFF/MasterScripts/10.0.0.1/proxy.service
[Unit]
Description=Kube Proxy
After=etcd_custom.service
After=docker_custom.service
Wants=etcd_custom.service
Wants=docker_custom.service

[Service]
ExecStart=/opt/bin/kube-proxy --etcd_servers=http://198.162.52.217:4001 --logtostderr=true

[Install]
